# services/resources/config/etc/cron.d/wbrules_checker
# m h dom mon dow user  command
* * * * *   root    /etc/wb-rules/services/resources/scripts/wbrules_checker.sh >> /etc/wb-rules/services/log/wbrules_checker.log

